# firecracker-images/config/runner.service
[Unit]
Description=GitHub Actions Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStartPre=/opt/setup-runner.sh
ExecStart=/usr/local/bin/start-runner.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
